# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── .NET API + Angular SPA ────────────────────────────────────────────────────
[program:dotnet]
command=dotnet /app/BrockLee.dll
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/supervisor/dotnet.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/dotnet-err.log
stderr_logfile_maxbytes=5MB
priority=20
environment=
    ASPNETCORE_ENVIRONMENT="Production",
    ASPNETCORE_URLS="http://0.0.0.0:5477",
    DOTNET_RUNNING_IN_CONTAINER="true"

; ── Python FastAPI ML sidecar ─────────────────────────────────────────────────
; start-period is long (90s) because:
;   - ML model training runs at startup
;   - 5 models × NPS + Index = 10 pipelines trained
;   - Kaggle dataset blend = up to 11,698 rows
;   - Gradient Boosting 200 estimators takes ~10–20s
; .NET uses retry+circuit-breaker so it handles Python being slow to start

[program:python]
command=/opt/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info
directory=/python
autostart=true
autorestart=true
startretries=5
startsecs=90
stdout_logfile=/var/log/supervisor/python.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/python-err.log
stderr_logfile_maxbytes=5MB
priority=10
environment=
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1",
    PYTHONPATH="/python"

; ── Log directory setup ───────────────────────────────────────────────────────
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock